# service/openrc/netinfo_sender_s3
#!/sbin/openrc-run
# https://github.com/OpenRC/openrc/blob/master/service-script-guide.md
# env
export N_ENDPOINT="xxx"
export N_ACCESS_KEY_ID="xxx"
export N_SECRET_ACCESS_KEY="xxx"
export N_BUCKET="xxx"
export N_OBJECT_PATH="xxx"
export N_ENCRYPTION_KEY="xxx"
export N_INTERVAL="5m"

command="/usr/local/bin/netinfo"
command_args="send s3 -endpoint=${N_ENDPOINT} -path_style -access_key_id=${N_ACCESS_KEY_ID} -secret_access_key=${N_SECRET_ACCESS_KEY} -bucket=${N_BUCKET} -object_path=${N_OBJECT_PATH} -encryption_key=${N_ENCRYPTION_KEY} -interval=${N_INTERVAL}"
command_background=true
command_user="root"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
  need net
  use dns logger netmount
  after firewall
}
