# service/openrc/netinfo_sender_webdav
#!/sbin/openrc-run
# https://github.com/OpenRC/openrc/blob/master/service-script-guide.md
# env
export N_ENDPOINT="xxx"
export N_USERNAME="xxx"
export N_PASSWORD="xxx"
export N_FILEPATH="xxx"
export N_ENCRYPTION_KEY="xxx"
export N_INTERVAL="5m"

command="/usr/local/bin/netinfo"
command_args="send webdav -endpoint=${N_ENDPOINT} -username=${N_USERNAME} -password=${N_PASSWORD} -filepath=${N_FILEPATH} -encryption_key=${N_ENCRYPTION_KEY} -interval=${N_INTERVAL}"
command_background=true
command_user="root"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
  need net
  use dns logger netmount
  after firewall
}
